# supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
directory=/
command=./emotes/start.sh

[program:celery]
directory=/
command=celery -A emotes.wsgi.celery -c 1 --loglevel=DEBUG worker
